# systemd/tmux-resurrect-save.service
[Unit]
Description=Save tmux sessions before shutdown
DefaultDependencies=no
Before=shutdown.target reboot.target halt.target
After=tmux-server.service
Requires=tmux-server.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/true
ExecStop=/bin/bash /home/jclee/.tmux/systemd/tmux-resurrect-save.sh
TimeoutStopSec=30

[Install]
WantedBy=default.target
