# systemd/tmux-slack-bridge.service
[Unit]
Description=tmux Slack Bridge
After=network-online.target tmux-server.service
Wants=network-online.target
PartOf=tmux-server.service

[Service]
Type=simple
ExecStart=%h/.tmux/bin/tmux-slack-bridge-start
WorkingDirectory=%h/.tmux/slack/tmux-bridge
EnvironmentFile=%h/.tmux/slack/tmux-bridge/.env
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
